# com.aiBrief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aiBrief.dailyUpdate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/baoying/Desktop/AI简报/daily-update.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/baoying/Desktop/AI简报/logs/update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/baoying/Desktop/AI简报/logs/error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/baoying/Desktop/AI简报</string>
</dict>
</plist>
